# scripts/queue/queue.service
[Unit]
Description=Coursero Queue Processor Service
After=network.target mysql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/coursero/scripts/queue
ExecStart=/usr/bin/php /var/www/coursero/scripts/queue/queue_processor.php
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=coursero-queue

[Install]
WantedBy=multi-user.target
